# install/etc/services.d/php-fpm/run
#!/usr/bin/with-contenv bash

mkdir -p /var/log/php
chown -R nginx:nginx /var/log/php
exec php-fpm7 -R -c /etc/php7/php-fpm.conf