# install/etc/services.d/php-fpm/run
#!/usr/bin/with-contenv bash

mkdir -p /var/log/php
chown -R www-data:www-data /var/log/php
exec php-fpm5 -R -c /etc/php5/php-fpm.conf